# src/supervisord.ini
[supervisord]
nodaemon=true

[program:notification]
command=python notification_service/telegram/broker.py
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stopsignal=QUIT
environment=PYTHONUNBUFFERED=1

[program:celery]
command=celery -A celery_app worker --loglevel=info
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stopsignal=QUIT
environment=PYTHONUNBUFFERED=1
